# crontab
0 *	* * *	root	/home/xp/start_proxy.sh
@reboot root /home/xp/start_proxy.sh
